# Scripts/restart_flask.service
[Unit]
Description=Restart Flask App

[Service]
Type=oneshot
ExecStart=/home/nfrd/Server/restart_flask.sh
